# environments_harbor/ssh_connection_refused_after_reboot_hard/environment/etc/systemd/system/ssh.service.d/override.conf
[Unit]
After=network-online.target network-ready.target
Wants=network-online.target
Before=multi-user.target

[Service]
Restart=on-failure
RestartSec=5s