# environments_harbor/systemd_timer_not_triggering_hard/environment/etc/systemd/system/database-backup.timer
[Unit]
Description=Database Backup Timer
Documentation=man:systemd.timer

[Timer]
OnCalendar=daily 02:00
Unit=database-backup.service

[Install]
WantedBy=timers.target