# environments_harbor/systemd_timer_not_triggering_hard/environment/etc/systemd/system/database-backup.service
[Unit]
Description=Database Backup Service
Documentation=man:systemd.service(5)
After=network.target postgresql.service

[Service]
Type=simple
User=backup
ExecStart=/usr/local/bin/backup-databse.sh
WorkingDirectory=/var/backups/database
StandardOutput=journal
StandardError=journal
TimeoutStartSec=300