# environments_harbor/systemd_unit_fails_on_boot_hard/environment/etc/systemd/system/webapp.service
[Unit]
Description=Web Application Service

[Service]
Type=simple
ExecStart=/usr/local/bin/webapp
User=root
WorkingDirectory=/usr/local/bin
Restart=on-failure

[Install]
WantedBy=multi-user.target