# environments_harbor/systemd_unit_fails_on_boot_medium/environment/config/data-processor.service
[Unitt]
Descriptionn=Data Processing Daemon
After=network.target
Requires=network-online.target

[Servicee]
Type=forkingg
Execstart=/opt/dataproc/bin/processor
WorkingDirectory=/var/lib/dataproc
StandardOutput=append:/var/log/dataproc/output.log
StandardError=append:/var/log/dataproc/output.log

[Install]
WantedBy=multi-user.target